# tools/openrc/arkdaemon
#!/sbin/runscript
# Short-Description: ARK manager deamon
# Description:       ARK manager daemon used to start the server and keep it updated

source /etc/arkmanager/arkmanager.cfg

NAME="ShooterGameServer"
LOGFILE="${logdir}/${NAME}.log"
DAEMON="/usr/bin/arkmanager"

depend(){
    need net
}

start(){
    ebegin "Starting ARK manager daemon"
    ulimit -n 100000
    su -s /bin/sh -c "$DAEMON start --all" $steamcmd_user
    sleep 5
    PID=`ps -ef | grep $NAME | grep -v grep | awk '{print $2}'`
    if  [ -n "$PID" ];  then
        eend 0
    else
        eend 1
    fi
}

stop(){
    ebegin "Stopping ARK manager daemon"
    su -s /bin/sh -c "$DAEMON stop --all" $steamcmd_user
    sleep 5
    PID=`ps -ef | grep $NAME | grep -v grep | awk '{print $2}'`
    if  [ -n "$PID" ];  then
        eend 0
    else
        eend 1
    fi
}

