# tools/openrc/arkdaemon
#!/sbin/runscript
# Short-Description: ARK manager deamon
# Description:       ARK manager daemon used to start the server and keep it updated

source /etc/arkmanager/arkmanager.cfg

NAME="ShooterGameServer"
LOGFILE="${logdir}/${NAME}.log"
DAEMON="/usr/bin/arkmanager"

depend(){
    need net
}

start(){
    INSTANCE="${RC_SVCNAME#*.}"
    ebegin "Starting ARK manager daemon"
    ulimit -n 100000
    if [ "$INSTANCE" != "$RC_SVCNAME" ]; then
        "$DAEMON" start "@${INSTANCE}"
    else
        "$DAEMON" start "@all"
    fi
    eend $?
}

stop(){
    INSTANCE="${RC_SVCNAME#*.}"
    ebegin "Stopping ARK manager daemon"
    if [ "$INSTANCE" != "$RC_SVCNAME" ]; then
        "$DAEMON" stop "@${INSTANCE}"
    else
        "$DAEMON" stop "@all"
    fi
    eend $?
}

