# tools/upstart/arkmanager.conf
start on runlevel [345]
stop on runlevel [!345]

env DAEMON="/usr/bin/arkmanager"

script
  for service in $(${DAEMON} list-instances --brief); do
    start arkmanager-instance service=${service}
  done
end script
